# expense-bot.service
[Unit]
Description=Telegram Expense Tracking Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/expense-bot
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/expense-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
